# deploy/quadlets/bspds-valkey.container
[Unit]
Description=BSPDS valkey cache
[Container]
ContainerName=bspds-valkey
Image=docker.io/valkey/valkey:9-alpine
Pod=bspds.pod
Volume=/srv/bspds/valkey:/data:Z
Exec=valkey-server --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru
HealthCmd=valkey-cli ping
HealthInterval=10s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=5s
[Service]
Restart=always
RestartSec=10
[Install]
WantedBy=default.target
