# deploy/quadlets/tranquil-pds-minio.container
[Unit]
Description=Tranquil PDS minio object storage
[Container]
ContainerName=tranquil-pds-minio
Image=cgr.dev/chainguard/minio:latest
Pod=tranquil-pds.pod
Environment=MINIO_ROOT_USER=minioadmin
Secret=tranquil-pds-minio-password,type=env,target=MINIO_ROOT_PASSWORD
Volume=/srv/tranquil-pds/minio:/data:Z
Exec=server /data --console-address :9001
[Service]
Restart=always
RestartSec=10
[Install]
WantedBy=default.target
